# deploy/systemd/stock-ts-daily-research.timer
[Unit]
Description=Refresh StockTS market and opportunity research

[Timer]
OnCalendar=*-*-* 07:20:00
OnCalendar=*-*-* 12:10:00
OnCalendar=*-*-* 18:30:00
Persistent=true
AccuracySec=5m
Unit=stock-ts-daily-research.service

[Install]
WantedBy=timers.target
