# deploy/systemd/stock-ts-morning-email.timer
[Unit]
Description=Send StockTS morning recap email every morning

[Timer]
OnCalendar=*-*-* 08:30:00
Persistent=true
AccuracySec=5m
Unit=stock-ts-morning-email.service

[Install]
WantedBy=timers.target
