# deploy/systemd/stock-ts-daily-analysis.timer
[Unit]
Description=Run StockTS data pipeline at market research checkpoints

[Timer]
OnCalendar=*-*-* 00:00:00
OnCalendar=*-*-* 09:00:00
OnCalendar=*-*-* 13:00:00
Persistent=true
AccuracySec=5m
Unit=stock-ts-daily-analysis.service

[Install]
WantedBy=timers.target
