# deploy/systemd/stock-ts-daily-analysis.service
[Unit]
Description=StockTS daily full data pipeline and report
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=admin
WorkingDirectory=/opt/stock-ts
Environment=PYTHONPATH=/opt/stock-ts/src
TimeoutStartSec=3600
ExecStart=/opt/stock-ts/.venv/bin/python scripts/run_daily_pipeline.py --python /opt/stock-ts/.venv/bin/python --tdx-bridge-python python3.11 --snapshot data/imports/tdx_snapshots.json --holdings data/portfolio/holdings.csv --provider tdx-snapshot --candidate-limit 300 --enrich-limit 50 --external-enrich-timeout 300 --announcement-limit 5 --output-dir reports/daily --html-dir reports/html --announcement-dir reports/announcements
